# bin/no_turbo.service
[Unit]
Description=Disable Intel Turbo Boost
After=multi-user.target
[Service]
Type=simple
ExecStart=echo 1 > /sys/devices/system/cpu/intel_pstate/no_turbo
Restart=always
RestartSec=20
StandardOutput=append:/var/log/powercontrol.log
StandardError=append:/var/log/powercontrol.log
[Install]
WantedBy=multi-user.target
