# bin/no_turbo.service
[Unit]
Description=Disable Intel Turbo Boost
After=multi-user.target
[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c 'echo 1 > /sys/devices/system/cpu/intel_pstate/no_turbo'
RemainAfterExit=yes
[Install]
WantedBy=multi-user.target
