# bin/powercontrol.service
[Unit]
Description=PowerControl
After=multi-user.target
[Service]
Type=simple
ExecStart=/usr/local/bin/powercontrol __powercontrol_monitor__
StandardOutput=append:/var/log/powercontrol.log
StandardError=append:/var/log/powercontrol.log
[Install]
WantedBy=multi-user.target
